# bin/commands/freebsd/countly_api
#!/bin/sh
#
# PROVIDE: countly_api
# REQUIRE: LOGIN DAEMON NETWORKING mountcritlocal
# KEYWORD: shutdown
#

. /etc/rc.subr

name=countly_api
rcvar=${name}_enable
pidfile=/var/run/${name}.pid

command="/usr/sbin/daemon"
command_args="-o ../../log/countly-api.log -P ${pidfile} -u countly -r /usr/local/bin/node ../../api/api.js"

load_rc_config ${name}

: ${countly_api_enable="NO"}
: ${countly_api_chdir=/usr/home/countly/countly-server/bin/config}

run_rc_command "$1"
