# bin/commands/freebsd/countly_dashboard
#!/bin/sh
#
# PROVIDE: countly_dashboard
# REQUIRE: LOGIN DAEMON NETWORKING mountcritlocal
# KEYWORD: shutdown
#

. /etc/rc.subr

name=countly_dashboard
rcvar=${name}_enable
pidfile=/var/run/${name}.pid

command="/usr/sbin/daemon"
command_args="-o /var/log/${name}.log -P ${pidfile} -u countly -r /usr/local/bin/node frontend/express/app.js"

load_rc_config ${name}

: ${countly_dashboard_enable="NO"}
: ${countly_dashboard_chdir=/usr/home/countly/countly-server}

run_rc_command "$1"
