# bin/commands/systemd/countly.service
[Unit]
Description=countly-supervisor
Documentation=http://count.ly

[Install]
WantedBy=multi-user.target

[Service]
LimitNOFILE=10000
Type=forking
Restart=on-failure
