# bin/commands/systemd/mongod.service
[Unit]
Description=An object/document-oriented database
Documentation=man:mongod(1)

[Install]
WantedBy=multi-user.target

[Service]
User=mongodb
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitNOFILE=64000
LimitRSS=infinity
LimitNPROC=32000
TimeoutStopSec=300
Restart=on-failure