# bin/commands/upstart/countly-supervisor.conf
description	"countly-supervisor"

limit nofile 10000 10000
limit nproc 10000 10000

start on runlevel [2345]
stop on runlevel [!2345]

respawn
